# arch/etc/systemd/system/snapper-timeline.timer
[Unit]
Description=Clean btrfs snapshots via snapper

[Timer]
OnCalendar=
OnCalendar=*:0/15

[Install]
WantedBy=timers.target
